# usr/lib/systemd/system/caddy.service.d/10-wait-for-domain.conf
[Unit]
ConditionPathExists=/etc/caddy/environment

[Service]
EnvironmentFile=/etc/caddy/environment
