# usr/lib/systemd/system/var-lib-jellyseerr.mount
[Unit]
Description=Config mount
After=network.target
Wants=network-online.target

[Mount]
What=unas-pro.jdmarble.net:/var/nfs/shared/jellyseerr_config
Where=/var/lib/jellyseerr
Type=nfs
Options=auto

[Install]
WantedBy=multi-user.target
