# playbooks/templates/var-mnt-media_nfs.mount.j2
[Mount]
What={{ hostvars['file-server']['ansible_host'] }}:{{ hostvars['file-server']['media_path'] }}
Type=nfs
Options=auto
